# archetypes/basic/src/main/resources/archetype-resources/package/__serviceName__.service
[Unit]
Description=$serviceName
Requires=network.service

[Service]
ExecStart=/opt/$serviceName/${serviceName}.jar
WorkingDirectory=/opt/$serviceName
Restart=on-failure

[Install]
WantedBy=multi-user.target
